# bazarr/rootfs/etc/s6-overlay/s6-rc.d/bazarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Bazarr
# Runs Bazarr
# ==============================================================================
bashio::log.info 'Starting Bazarr...'

INGRESS_PORT=$(bashio::addon.ingress_port)
INGRESS_ENTRY=$(bashio::addon.ingress_entry)

bashio::log.info "Ingress port: ${INGRESS_PORT}"
bashio::log.info "Ingress entry: ${INGRESS_ENTRY}"

# Configure Bazarr base URL for ingress path in config.yaml
CFG_PATH="/config/config/config.yaml"

if [ -f "${CFG_PATH}" ]; then
    # Replace base_url only on the first occurrence (in the general section)
    sed -i "0,/  base_url:.*/s|  base_url:.*|  base_url: ${INGRESS_ENTRY}|" "${CFG_PATH}"
    bashio::log.info "Updated base_url to: ${INGRESS_ENTRY}"
else
    bashio::log.warning "config.yaml not found, Bazarr will create it on first run"
fi

python3 /opt/bin/bazarr.py \
    --no-update \
    --config /config \
    --port "${INGRESS_PORT}"
